# deploy_tools/gunicorn-supervisor.template.conf
[program:SITENAME]
user=USERNAME
directory=/home/USERNAME/sites/SITENAME/source
command=/home/USERNAME/sites/SITENAME/virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket --capture-output --access-logfile ../access.log --error-logfile ../error.log --timeout=60 superlists.wsgi:application
